# config/ansible/templates/gunicorn.service.j2
[Unit]
Description=Gunicorn instance for portfolio app
After=network.target

[Service]
User=ec2-user
Group=nginx
WorkingDirectory=/opt/portfolio
Environment="PATH=/opt/portfolio/venv/bin"
ExecStart=/opt/portfolio/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 app:app
Restart=always

[Install]
WantedBy=multi-user.target
